# overlay/etc/rc.local
#!/bin/sh

Xorg -configure
xinit /etc/xtasks.local

